# odhom.conf
description     "Odhom Server"

start on runlevel [2345]
stop on runlevel [016]

chdir /opt/odhom

respawn limit 10 5

script
    exec /opt/odhom/odhom-server.sh
end script

